# jdk/1.8/etc/profile.d/java.sh
#!/bin/bash

. /etc/profile.d/add-docker-env.sh
